# config/ethnos-api.service
[Unit]
Description=Ethnos Academic API Service
Documentation=https://github.com/ethnos-app/api
After=network.target mysql.service redis.service sphinxsearch.service
Wants=mysql.service redis.service sphinxsearch.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=server
Group=server
WorkingDirectory=/home/server/api
Environment=NODE_ENV=production
Environment=PATH=/home/server/.nvm/versions/node/v22.18.0/bin:/usr/local/bin:/usr/bin:/bin
ExecStartPre=/bin/bash -c 'cd /home/server/api && /bin/bash ./server.sh cleanup'
ExecStart=/bin/bash -c 'cd /home/server/api && /home/server/.nvm/versions/node/v22.18.0/bin/npm start'
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/bash -c 'cd /home/server/api && ./server.sh stop'
Restart=always
RestartSec=10
StandardOutput=append:/home/server/api/logs/systemd.log
StandardError=append:/home/server/api/logs/systemd-error.log
SyslogIdentifier=ethnos-api

# Resource limits
LimitNOFILE=65536
MemoryMax=2G
CPUQuota=200%

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/server/api/logs /home/server/api/uploads /tmp
ProtectHome=read-only
ProtectKernelTunables=true
ProtectControlGroups=true

# Health check
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
